# Language/Batch/Bash/Project/service_jar_upstart.conf
# 将jar包部署为Linux服务，并保证服务异常挂掉后会自动重启，适用于CentOS6操作系统。

author       "Michael Brucelin <michaelbrucelin@gmail.com>"
description  "Test Service"

start on runlevel [2345]
stop on runlevel [016]

respawn
respawn limit unlimited

env upstart_name=jar-service
env stdout_log=/var/log/jar-service.log
env stderr_log=/var/log/jar-service.err
# env pid_file=/var/run/jar-service.pid  # 不需要管理pid文件了

pre-start script
    echo "[`date`] jar-service Starting" | ts >>"$stdout_log"
end script

script
    # exec /usr/bin/java -jar /etc/jar-service/jar-service.jar | ts >>"$stdout_log" 2>>"$stderr_log"  # 使用ts后无法将stdout与stderr分开
    exec /usr/bin/java -jar /etc/jar-service/jar-service.jar 2>&1 | ts >>"$stdout_log"
end script

pre-stop script
    echo "[`date`] jar-service Stopping" | ts >>"$stdout_log"
end script