# Language/Batch/Bash/Project/service_jar_upstart.conf
# 将jar包部署为Linux服务，并保证服务异常挂掉后会自动重启，适用于CentOS6操作系统。

author       "Michael Brucelin <michaelbrucelin@gmail.com>"
description  "Test Service"

start on runlevel [2345]
stop on runlevel [016]

respawn
respawn limit unlimited
# respawn limit 10 5  # 设置服务异常停止后重启次数及间隔时间

script
    exec /usr/bin/java -jar /etc/jar-service/jar-service.jar
end script

pre-start script
    echo "[`date`] Node Test Starting" >> /var/log/nodetest.log
end script

pre-stop script
    echo "[`date`] Node Test Stopping" >> /var/log/nodetest.log
end script