# Language/Batch/Bash/Template/service_upstart1.conf
# 模板来自:https://gist.github.com/rezen/9792753

author       "Michael BruceLin <michaelbrucelin@gmail.com>"
description  "Update Example Config"

# Set variables
env upstart_name=template
env log_file=/var/log/upstart-template.txt
env pid_file=/var/run/template.pid

# If you want to pass args to upstart 
# e.g. start template src=/var/www/test
instance $src

# Start ....
start on startup
stop on shutdown

# Give up if restart occurs 10 times in 90 seconds.
respawn limit 10 90

# Before it starts ....
pre-start script
    PID=`status $upstart_name | egrep -oi '([0-9]+)$' | head -n1`
    echo $PID > $pid_file
    logger "$PID pre-starting .... "
    mkdir -p /var/www/html/upstart
end script

# The main thing that loops forever
script
    logger "the logfile variable = $log_file"
    sleep 1000
    logger "script stopped sleeping ... "
end script

# Before it all dies ...
pre-stop script
    logger "$PID pre-stopping .... "
    rm -rf /var/www/html/upstart
    rm -f $pid_file
end script